# systemd/layer-shell.service
[Unit]
Description=Layer Shell

[Service]
ExecStart=layer-shell
WorkingDirectory=%h
Restart=on-failure
RestartSec=5s
